# cookbooks/docker_netauth/files/default/rc.local
#!/bin/sh -e

# ipset restore < /etc/ipset
# iptables-restore /etc/iptables-rules

